# profile.d/stow.sh
# shellcheck shell=sh

#
# Shrc Stow

# https://www.gnu.org/software/stow/manual/stow.html
#
export STOW_DIR=""

