# etc/profile.d/config.sh
# shellcheck shell=sh

#
# Shrc Stow

# Stow Directory
# https://www.gnu.org/software/stow/manual/stow.html
#
export STOW_DIR

export DIRENV_CONFIG="${RC_CONFIG}/direnv"
export INPUTRC="${RC_CONFIG}/readline/inputrc"
export STARSHIP_CONFIG="${RC_CONFIG}/config.toml"
